# src/main/debian/etc/systemd/user/x11vnc@.service
[Unit]
Description=x11vnc Service

[Service]
Type=simple
ExecStart=/usr/sbin/systemd-x11vnc
StandardError=journal
StandardInput=socket
StandardOutput=socket
